# metalfan.service
[Unit]
Description=Fan control script

[Service]
Type=simple
ExecStart=/usr/bin/metalfan

[Install]
WantedBy=multi-user.target
